# etc/systemd/system/kino24-web.socket
[Unit]
Description=Socket for kino24 web service

[Socket]
ListenStream=127.0.0.1:10090

[Install]
WantedBy=sockets.target
